# ops/dev-tools-nav-deploy.timer
[Unit]
Description=Check dev-tools-nav main for a deployable SHA every ten minutes

[Timer]
OnBootSec=2min
OnCalendar=*:0/10
Persistent=true
Unit=dev-tools-nav-deploy.service

[Install]
WantedBy=timers.target
